# ng-ssl/crontab
# m	h	d	M	w	command
  0	4	*	*	*	certbot certonly --nginx -n -d ${SSL_DOMAIN}

