# mongo-dump/crontab
# m	h	d	M	w	command
0	4	*	*	* /usr/local/bin/mongo-dump
#
